# debian/postinst
#!/bin/sh
# postinst script for #PACKAGE#
#
# see: dh_installdeb(1)

#set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

configure() {
	# Divert lightdm configurations
	if ! [ -e "/etc/lightdm/lightdm.conf.distrib" ]; then
		dpkg-divert --rename --package semplice-openbox-default-settings --divert "/etc/lightdm/lightdm.conf.distrib" --add "/etc/lightdm/lightdm.conf"
		cp /usr/share/semplice-default-settings/lightdm/lightdm.conf /etc/lightdm
	fi
	if ! [ -e "/etc/lightdm/lightdm-gtk-greeter.conf.distrib" ]; then
		dpkg-divert --rename --package semplice-openbox-default-settings --divert "/etc/lightdm/lightdm-gtk-greeter.conf.distrib" --add "/etc/lightdm/lightdm-gtk-greeter.conf"
		cp /usr/share/semplice-default-settings/lightdm/lightdm-gtk-greeter.conf /etc/lightdm
	fi
	#invoke-rc.d lightdm reload
	
	# Divert chrome configurations
	if ! [ -e "/etc/chromium/master_preferences.distrib" ]; then
		dpkg-divert --rename --package semplice-openbox-default-settings --divert "/etc/chromium/master_preferences.distrib" --add "/etc/chromium/master_preferences"
		ln -s /usr/share/semplice-default-settings/chromium/master_preferences /etc/chromium/master_preferences
	fi
	
	# Update gconf defaults
	if [ "`which update-gconf-defaults`" ]; then
	 update-gconf-defaults
	fi
}

case "$1" in
    configure|semplice-live-configure)
	configure;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
